# config/systemd/user/castget.service
[Unit]
Description=Script to update newsbeuter

[Service]
Type=simple
ExecStart=/usr/bin/castget -1

[Install]
WantedBy=default.target
