# config/systemd/user/dwl-session.target
[Unit]
Description=dwl compositor session
Documentation=man:systemd.special
BindsTo=dwl.service
After=dwl.service
