# config/systemd/user/dwl.service
[Unit]
Description=dwl compositor
Wants=graphical-session.target dbus.service seatd.service
After=graphical-session.target dbus.service seatd.service
Wants=dwl-session.target

[Service]
Type=simple

# Canonical runtime + bus (so dbus-update-activation-environment talks to the right bus)
Environment=XDG_RUNTIME_DIR=%t
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=%t/bus

# wlroots / seatd
Environment=LIBSEAT_BACKEND=seatd
Environment=SEATD_SOCK=/run/seatd.sock
# (Optional but robust) force DRM so it never tries the Wayland backend
# Environment=WLR_BACKENDS=drm
Environment=WLR_LOG=info
Environment=XDG_SESSION_TYPE=wayland
Environment=XDG_CURRENT_DESKTOP=dwl

# Ensure compositor starts with no inherited client-display vars
# Environment=DISPLAY=
# Environment=WAYLAND_DISPLAY=
# UnsetEnvironment=DISPLAY WAYLAND_DISPLAY DBUS_SESSION_BUS_ADDRESS
UnsetEnvironment=DISPLAY

# 1) Nuke any stale wayland sockets from previous runs (safe: only your runtime dir)
ExecStartPre=/usr/bin/find %t -maxdepth 1 -type s -name 'wayland-*' -user %u -delete

# 2) Make sure the user manager has a clean env
ExecStartPre=/usr/bin/systemctl --user unset-environment DISPLAY WAYLAND_DISPLAY DBUS_SESSION_BUS_ADDRESS

# 3) Publish correct env to user-systemd (DBUS pinned to %t/bus above)
ExecStartPre=/usr/bin/dbus-update-activation-environment --systemd \
  XDG_RUNTIME_DIR DBUS_SESSION_BUS_ADDRESS XDG_SESSION_TYPE XDG_CURRENT_DESKTOP

# Launch dwl (with your startup script if you like)
ExecStart=/usr/bin/dwl -s %h/.local/bin/dwl-startup.sh
# ExecStart=/bin/sh -lc '/usr/local/bin/slstatus -s | exec /home/go/src/codeberg.org/dwl/dwl/dwl -s /home/go/.local/bin/dwl-startup.sh'

# ExecStartPost=systemctl --user start dwl-session.target

# ExecStart=/usr/local/bin/dwl
# 

# Keep things tidy on exit so next start is clean
ExecStopPost=/usr/bin/systemctl --user unset-environment DISPLAY WAYLAND_DISPLAY
# ExecStopPost=/usr/bin/systemctl --user stop dwl-session.target
ExecStopPost=/usr/bin/systemctl --user stop dwl-session.target
ExecStopPost=/usr/bin/systemctl --user unset-environment DISPLAY WAYLAND_DISPLAY

Restart=on-failure
RestartSec=1
ExecStop=/usr/bin/pkill -TERM dwl
TimeoutStopSec=5

[Install]
WantedBy=graphical-session.target
