# config/systemd/user/swayidle.path
# ~/.config/systemd/user/swayidle.path
[Unit]
Description=Start swayidle when Wayland socket appears
PartOf=dwl.service
After=dwl.service

[Path]
PathExistsGlob=%t/wayland-*

[Install]
WantedBy=graphical-session.target

