# config/systemd/user/task-sync.service
[Unit]
Description=Run task sync every few minutes

[Service]
Type=simple
ExecStart=/usr/bin/task sync
NoNewPrivileges=true
