# config/systemd/user/waybar.path
[Unit]
Description=Start Waybar when Wayland socket appears
PartOf=dwl.service
After=dwl.service

[Path]
PathExistsGlob=%t/wayland-*

[Install]
WantedBy=graphical-session.target
